# units/user-configdrive.path
[Unit]
Description=Watch for a cloud-config at /media/configdrive
Requires=media-configdrive.mount
After=media-configdrive.mount

[Path]
DirectoryNotEmpty=/media/configdrive
